# templates/deb/postinst.erb
#!/bin/bash
<%- get_services.each do |service| -%>
  # switch based on systemd vs systemv
  #
  <%- if @platform.servicetype == "systemd" -%>
if [ -z "$2" ]; then
  systemctl enable <%= service.name %>.service >/dev/null || :
else
  systemctl try-restart <%= service.name %>.service >/dev/null || :
fi
  <%- elsif @platform.servicetype == "sysv" -%>
if [ -x "<%= service.service_file %>" ]; then
  update-rc.d <%= service.name %> defaults > /dev/null

  if [ -n "$2" ]; then
    invoke-rc.d <%= service.name %> condrestart || true
  fi
fi
  <%- end -%>
<%- end -%>

<%- if @user -%>
# Add our user and group
<%= @platform.add_group(@user) %>
<%= @platform.add_user(@user) %>
<%- end -%>

# Set up any specific permissions needed...
<%- (get_directories + get_configfiles + get_files).select { |pathname| pathname.has_overrides? }.uniq.each do |file_or_directory| -%>
  <%= "chmod '#{file_or_directory.mode}' '#{file_or_directory.path}'" if file_or_directory.mode %>
  <%- if file_or_directory.owner -%>
    if getent passwd '<%= file_or_directory.owner %>' &> /dev/null; then
      chown '<%= file_or_directory.owner %>' '<%= file_or_directory.path %>'
    else
      echo "Error updating '<%= file_or_directory.path %>': user '<%= file_or_directory.owner %>' does not exist."
    fi
  <%- end -%>
  <%- if file_or_directory.group -%>
    if getent group '<%= file_or_directory.group %>' &> /dev/null; then
      chgrp '<%= file_or_directory.group %>' '<%= file_or_directory.path %>'
    else
      echo "Error updating '<%= file_or_directory.path %>': group '<%= file_or_directory.group %>' does not exist."
    fi
  <%- end -%>
<%- end -%>

# Run postinstall scripts on install if defined
<%- scripts = get_postinstall_actions("install") -%>
<%- unless scripts.empty? -%>
if [ $1 = configure ] && [ -z $2] ; then
  <%= scripts %>
fi
<%- end -%>

# Run postinstall scripts on upgrade if defined
<%- scripts = get_postinstall_actions("upgrade") -%>
<%- unless scripts.empty? -%>
if [ $1 = configure ] && [ -n $2] ; then
  <%= scripts %>
fi
<%- end -%>
