# templates/deb/postrm.erb
#!/bin/bash

# Run post-removal scripts on upgrade if defined
<%- scripts = get_postremove_actions("upgrade") -%>
<%- unless scripts.empty? -%>
if [ $1 = upgrade ] && [ -n $2] ; then
  <%= scripts %>
fi
<%- end -%>

# Run post-removal scripts on removal if defined
<%- scripts = get_postremove_actions("removal") -%>
<%- unless scripts.empty? -%>
if [ $1 = remove ] && [ -z $2] ; then
  <%= scripts %>
fi
<%- end -%>

<%- get_services.each do |service| -%>
  # switch based on systemd vs systemv
  #
  <%- if @platform.servicetype == "systemd" -%>
systemctl daemon-reload >/dev/null 2>&1 || :
  <%- elsif @platform.servicetype == "sysv" -%>
if [ "$1" = "purge" ] ; then
  update-rc.d <%= service.name %> remove > /dev/null
fi
  <%- end -%>
<%- end -%>
