# templates/deb/prerm.erb
#!/bin/bash

# Run pre-removal scripts on upgrade if defined
<%- scripts = get_preremove_actions("upgrade") -%>
<%- unless scripts.empty? -%>
if [ $1 = upgrade ] && [ -n $2] ; then
  <%= scripts %>
fi
<%- end -%>

# Run pre-removal scripts on removal if defined
<%- scripts = get_preremove_actions("removal") -%>
<%- unless scripts.empty? -%>
if [ $1 = remove ] && [ -z $2] ; then
  <%= scripts %>
fi
<%- end -%>

<%- get_services.each do |service| -%>
  # switch based on systemd vs systemv
  #
  <%- if @platform.servicetype == "systemd" -%>
if [ "$1" = remove ]; then
  systemctl --no-reload disable <%= service.name %>.service > /dev/null 2>&1 || :
  systemctl stop <%= service.name %>.service > /dev/null 2>&1 || :
fi

  <%- elsif @platform.servicetype == "sysv" -%>
if [ -x "<%= service.service_file %>" ] && [ "$1" = remove ]; then
  invoke-rc.d <%= service.name %> stop || true
fi
  <%- end -%>
<%- end -%>
