# resources/deb/postrm.erb
#!/bin/bash

# Run post-removal scripts on upgrade if defined
if [ $1 = upgrade ] && [ -n $2 ] ; then
  <%= get_postremove_actions("upgrade") %>
fi

# Run post-removal scripts on removal if defined
if [ $1 = remove ] && [ -z $2 ] ; then
  <%= get_postremove_actions("removal") %>
fi

<%- get_services.each do |service| -%>
  # switch based on systemd vs systemv
  #
if [ -d '/run/systemd/system' ] ; then
  systemctl daemon-reload >/dev/null 2>&1 || :
else
  if [ "$1" = "purge" ] ; then
    update-rc.d <%= service.name %> remove > /dev/null
  fi
fi
<%- end -%>
